# etc/udev/rules.d/monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", RUN+="/usr/local/bin/hotplug-monitor.sh"
